# conf/daemons/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:mysqld]
command=/bin/mysqld_safe
stdout_events_enabled=true
stderr_events_enabled=true

[program:httpd]
command=/usr/sbin/httpd
stdout_events_enabled=true
stderr_events_enabled=true

[program:sshd]
command=/usr/sbin/sshd
stdout_events_enabled=true
stderr_events_enabled=true

[program:crond]
command=/usr/sbin/crond
stdout_events_enabled=true
stderr_events_enabled=true

[program:mailcatcher]
command=/usr/local/bin/mailcatcher --ip=0.0.0.0 --http-port=81
stdout_events_enabled=true
stderr_events_enabled=true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
